# dashboard/dspark-live-dashboard.service.in
# SPDX-License-Identifier: MIT
[Unit]
Description=DSpark live vLLM dashboard
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
User=@USER@
Group=@GROUP@
WorkingDirectory=@REPO_ROOT@
Environment=PYTHONDONTWRITEBYTECODE=1
ExecStart=@REPO_ROOT@/dashboard/run-dashboard.sh
Restart=always
RestartSec=3
PrivateTmp=true

[Install]
WantedBy=multi-user.target
